# supervisord.ini
[supervisord]
nodaemon=true

[program:X11]
command=/usr/bin/Xvfb :0 -screen 0 1024x768x24
autorestart=true

[program:x11vnc]
command=/usr/bin/x11vnc -display :0 -ncache 10
autorestart=true

[program:openbox]
user=user
environment=DISPLAY=":0",HOME="/home/user"
command=/usr/bin/openbox-session
autorestart=true

[program:noVNC]
command=/opt/noVNC/utils/novnc_proxy --vnc localhost:5900 --listen 8083
workingdir=/opt/noVNC
autorestart=true