# scripts/supervisord.ini
[supervisord]
nodaemon=true

[eventlistener:log]
priority=0
buffer_size=500
command=/usr/bin/python /scripts/log.py
events=PROCESS_LOG
stderr_logfile=/dev/stdout
stderr_logfile_maxbytes=0

[program:X11]
priority=100
command=/usr/bin/Xvfb :0 -screen 0 1024x768x16
autorestart=true
stdout_events_enabled=true
stderr_events_enabled=true

[program:tigerVNC]
priority=200
command=/scripts/x0vncserver
autorestart=true
stdout_events_enabled=true
stderr_events_enabled=true

[program:noVNC]
priority=300
command=/scripts/noVNC
autorestart=true
stdout_events_enabled=true
stderr_events_enabled=true

[program:openbox]
priority=400
environment=DISPLAY=":0"
command=/usr/bin/sudo -iu user /usr/bin/openbox-session
autorestart=true
stdout_events_enabled=true
stderr_events_enabled=true
