# scripts/supervisord.ini
[supervisord]
nodaemon=true

[program:X11]
priority=100
command=/usr/bin/Xvfb :0 -screen 0 1024x768x16
autorestart=true

[program:tigerVNC]
priority=200
command=/scripts/x0vncserver
autorestart=true

[program:noVNC]
priority=300
user=daemon
command=/opt/noVNC/utils/novnc_proxy --vnc localhost:5900 --listen 8083
autorestart=true

[program:openbox]
priority=400
user=user
environment=DISPLAY=":0",HOME="/home/user",USER="user"
command=/usr/bin/openbox-session
autorestart=true
