# scripts/supervisord.ini
[supervisord]
nodaemon=true

[program:X11]
priority=100
command=/usr/bin/Xvfb :0 -screen 0 1024x768x16
autorestart=true

[program:tigerVNC]
priority=200
command=/scripts/x0vncserver
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:noVNC]
priority=300
command=/scripts/noVNC
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:openbox]
priority=400
environment=DISPLAY=":0"
command=/usr/bin/sudo -iu user /usr/bin/openbox-session
autorestart=true
