# scripts/supervisord.ini
[supervisord]
nodaemon=true

[program:X11]
command=/usr/bin/Xvfb :0 -screen 0 1280x800x24
autorestart=true

[program:x11vnc]
command=/scripts/x11vnc
autorestart=true

[program:openbox]
user=user
environment=DISPLAY=":0",HOME="/home/user",USER="user"
command=/usr/bin/openbox-session
autorestart=true

[program:noVNC]
user=daemon
command=/opt/noVNC/utils/novnc_proxy --vnc localhost:5900 --listen 8083
autorestart=true